# outdir/softwares/expr.cron
2 * * * * /usr/bin/flock -w 0 /home/$USER/sptestlock bash outdir/softwares/speedtest/runtestround.bash >lastoutput 2>&1
